# factory/usr/lib/systemd/system/the-tool@data.service.d/override.conf
[Unit]
Wants=run-mnt-ubuntu\x2data.mount
After=run-mnt-ubuntu\x2data.mount
After=the-tool@run.service
After=the-tool@install.service

