# factory/usr/lib/systemd/system/initrd-switch-root.service.d/core-override.conf
[Service]
ExecStart=
#Ensure that we do not perform stateful reexec (boot fails)
#Ensure that we ignore kernel cmdline init= parameter (dangerous)
ExecStart=/bin/systemctl --no-block switch-root /sysroot /sbin/init
