# factory/usr/lib/systemd/system/snapd.recovery-chooser-trigger.service.d/fake-clean-stop.conf
# TODO: We need to make `snap-bootstrap recovery-chooser-trigger` to
# survive the switch root. That means receive a signal so it chroots
# and chdirs as well as setting argv[0]'s first character to '@'.  For
# the moment we make kill a proper termination status so that the
# service does not go in "failed" state which prevents main boot to
# start it properly.

[Service]
SuccessExitStatus=SIGTERM 0
