# factory/usr/lib/systemd/system/initrd-mode@.target
[Unit]
OnFailure=emergency.target
OnFailureJobMode=replace-irreversibly
Wants=basic.target
After=basic.target
Wants=the-tool@%i.service
