# factory/usr/lib/systemd/system/modeenv-available.target
[Unit]
Description=Modeenv available

Wants=reload-disk.service
After=reload-disk.service
