# factory/usr/lib/systemd/system/snap-update-boot-loader-vars.service
[Unit]
ConditionKernelCommandLine=snapd_recovery_mode=run
Description=Update bootloader variables
OnFailure=emergency.target
OnFailureJobMode=replace-irreversibly
DefaultDependencies=no

# run-mnt-boot.mount is generated by /usr/lib/systemd/system-generators/snap-mounts
After=run-mnt-boot.mount

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/lib/snapd/snap-bootstrap update-boot-loader-vars
StandardOutput=journal+console
StandardError=journal+console
