# factory/usr/lib/systemd/system/snap-data-save-mounted.target
[Unit]
Description=Data mounted

Wants=reload-model.service
After=reload-model.service
