# factory/usr/lib/systemd/system/snap-compare-data-save.service
[Unit]
ConditionKernelCommandLine=snapd_recovery_mode=run
ConditionPathIsEncrypted=/run/mnt/data
Description=Compare data and save partitions
OnFailure=emergency.target
OnFailureJobMode=replace-irreversibly
DefaultDependencies=no

After=reload-model.service

# run-mnt-data.mount and run-mnt-ubuntu\x2dsave.mount are generated by /usr/lib/systemd/system-generators/snap-mounts
After=run-mnt-ubuntu\x2dsave.mount
After=run-mnt-data.mount

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/lib/snapd/snap-bootstrap compare-data-save
StandardOutput=journal+console
StandardError=journal+console
