# factory/usr/lib/systemd/system/sysroot.mount
[Unit]
Before=initrd-cleanup.target
[Mount]
What=/run/mnt/base
Type=none
Options=bind
