# factory/usr/lib/systemd/system/sysroot.mount
[Unit]
Before=initrd-cleanup.target
DefaultDependencies=no
After=the-tool.service
[Mount]
What=/run/mnt/base
Type=none
Options=bind
