# factory/usr/lib/systemd/system/sysroot-writable.mount
[Unit]
Before=initrd-cleanup.target
DefaultDependencies=no
After=run-mnt-data.mount
After=run-mnt-base.mount
After=sysroot.mount
After=the-tool.service
Requires=sysroot.mount
Requires=the-tool.service
Wants=populate-writable.service
[Mount]
What=/run/mnt/data
Type=none
Options=bind
