# factory/usr/lib/systemd/system/the-tool.service
[Unit]
ConditionKernelCommandLine=snapd_recovery_mode

[Service]
Type=oneshot
RemainAfterExt=true
ExecStart=/usr/lib/the-tool
