# factory/usr/lib/systemd/system/the-tool.service
[Unit]
ConditionKernelCommandLine=snapd_recovery_mode
After=run-mnt-ubuntu\x2dseed.mount
Requires=run-mnt-ubuntu\x2dseed.mount
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/lib/the-tool
