# factory/usr/lib/systemd/system/the-tool.service
[Service]
Type=oneshot
RemainAfterExt=true
ExecStart=/usr/lib/the-tool
