# factory/usr/lib/systemd/system/initrd-cleanup.service.d/core-override.conf
[Service]
ExecStart=
ExecStart=systemctl --no-block start initrd-switch-root.target
