# fips/usr/lib/systemd/system/snap-initramfs-mounts.service.d/fips-img.conf
[Unit]
After=systemd-modules-load.service

# This is for kernels that are not booted from an EFI environment. This probably
# means a FIT image on arm64.
ConditionFirmware=!uefi

[Service]

# Ensure tcrypt was loaded
ExecStartPre=modprobe tcrypt
# Then allow snap-bootstrap to mount things
# Self-test verify FIPS computed hmac
ExecStartPost=kcapi-hasher -n sha512hmac --check /run/mnt/kernel/.kernel.img.hmac
# On failure, emergency.target is triggered with 10 second delay
